# module/DoorBundle/Resources/bin/door.service
[Unit]
Description=door
After=network-online.target
Wants=network-online.target

[Service]
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python /usr/local/sbin/door.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
